# dotfiles/smart_plug/linux/smart-plug-daemon.service
[Unit]
Description=Smart Plug Daemon

[Service]
ExecStart=/home/biggs/.nix-profile/bin/speakerctl

[Install]
WantedBy=multi-user.target
