# dotfiles/smart_plug/linux/smart-plug-daemon.service
[Unit]
Description=Smart Plug Daemon

[Service]
ExecStart=/opt/speaker/speakerctl

[Install]
WantedBy=multi-user.target
