# dotfiles/nix/systemd-garbage-collection/nix-garbage-collection.timer
[Unit]
Description=Collect garbage in the Nix store

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
